# conf/gunicorn.conf
[program:gunicorn]
directory=/home/hde/fintest
command=/home/hde/fintest/scripts/start.sh
autorestart=true
autostart=true
redirect_stderr=true
user=hde